# fpms.service.in
[Unit]
Description=FPMS daemon responsible for OLED management

[Service]
ExecStart=python3 @datadir@/app/fpms.py
WorkingDirectory=@datadir@/app

[Install]
WantedBy=multi-user.target
