# fpms.service.in
[Unit]
Description=FPMS daemon responsible for OLED management

[Service]
ExecStart=python3 @datadir@/fpms/fpms.py
WorkingDirectory=@datadir@/fpms

[Install]
WantedBy=multi-user.target
