# debian/wlanpi-fpms.service
[Unit]
Description=Front Panel Menu System for the WLAN Pi OLED display

[Service]
ExecStart=/usr/sbin/fpms

[Install]
WantedBy=multi-user.target
